# scripts/systemd/fs-mission.service
# Reference unit -- the front-seat mission, started on demand by bb_init.sh
# (via `systemctl start fs-mission.service`). You maintain the real copy of
# this; it is included here so bb_init.sh's launch step is self-documenting.
#
# Install with:
#     sudo cp scripts/systemd/fs-mission.service /etc/systemd/system/
#     sudo systemctl daemon-reload
#   (Do NOT `enable` it -- bb_init.sh starts it only after the gates pass.)

[Unit]
Description=BlueBoat front-seat MOOS mission
# NO ordering on bb-init.service. bb_init.sh starts this unit explicitly via a
# synchronous `systemctl start`; an After=bb-init dep here deadlocks, because
# bb-init (Type=oneshot) stays "activating" until its script returns, while the
# script blocks waiting for this start job that systemd won't run until bb-init
# is active. fs-mission is on-demand only -- it must not be ordered after the
# thing that launches it.

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/moos-ivp-greece/missions/frontseat_greece
# A systemd unit has a bare PATH; launch_vehicle.sh needs pAntler, nsplug,
# mhash_gen, uMAC, and the repo's own binaries. List every dir they live in.
Environment=PATH=/home/pi/moos-ivp-blueboat/bin:/home/pi/moos-ivp/bin:/home/pi/moos-ivp/scripts:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# -a makes launch_vehicle.sh `exec pAntler` in the FOREGROUND (no uMAC); with the
# `exec` below, pAntler becomes the service's main process so systemd tracks it
# directly -- output to journald, restart-capable, no screen in the critical path.
# -nc skips the interactive confirm. (We tried wrapping this in `screen` for an
# attachable console; screen-under-systemd proved flaky, so it's reverted.)
#
# To watch or drive the running mission after logging in, WITHOUT touching the
# service lifecycle:
#   live output:       journalctl -fu fs-mission.service
#   interactive uMAC:  uMAC targs/targ_<vname>.moos   (e.g. targ_yip.moos, in the
#                      mission dir) -- uMAC just connects to the live MOOSDB, so
#                      you can open/close it freely.
ExecStart=/bin/bash -c 'exec ./launch_vehicle.sh -a -nc'
# Start-once semantics: bb_init.sh decides when to (re)start; do not loop here.
Restart=no

[Install]
WantedBy=multi-user.target
